# deploy/systemd/symphonyd-maintenance.timer
[Unit]
Description=Daily symphonyd SQLite backup and log pruning
Documentation=https://github.com/kulichevskiy/symphonyd

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m
Unit=symphonyd-maintenance.service

[Install]
WantedBy=timers.target
