# deploy/systemd/symphonyd-maintenance.service
[Unit]
Description=symphonyd SQLite backup and log pruning
Documentation=https://github.com/kulichevskiy/symphonyd
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=symphony
Group=symphony
WorkingDirectory=/opt/symphonyd
EnvironmentFile=-/opt/symphonyd/.env
Environment=HOME=/home/symphony
Environment=PATH=/opt/symphonyd/.venv/bin:/home/symphony/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=SYMPHONYD_CONFIG=/opt/symphonyd/config.yaml
Environment=SYMPHONYD_BACKUP_KEEP=7
Environment=SYMPHONYD_LOG_RETENTION_DAYS=14
ExecStartPre=/usr/bin/test -f /opt/symphonyd/config.yaml
ExecStartPre=/usr/bin/test -x /opt/symphonyd/.venv/bin/python
ExecStartPre=/usr/bin/test -x /usr/bin/sqlite3
ExecStart=/opt/symphonyd/.venv/bin/python /opt/symphonyd/deploy/systemd/symphonyd-maintenance.py
NoNewPrivileges=true
PrivateTmp=true
